# str2str/str2str.service
# str2str.service
# autogenerated by https://github.com/symsyak/rtk-base-station mgmt-cli

[Unit]
Description=rtk-base-station str2str
Wants=network-online.target
After=network-online.target
Conflicts=shutdown.target

StartLimitInterval=10s
StartLimitBurst=20

[Service]
Restart=always
RestartSec=1s
ExecStart=./opt/rtk-base-station/str2str/str2str.sh

[Install]
WantedBy=default.target