# sync-docker-compose/sync-docker-compose.timer
[Unit]
Description=timer sync.py

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target